# doc/thrivedevcenter-clean-sessions.timer
[Unit]
Description=Run ThriveDevCenter session clean daily

[Timer]
OnCalendar=*-*-* 10:38:05

[Install]
WantedBy=timers.target